# crontab
*/5 * * * *  python /app/retornaPontosDB.py 
